# hamweather/scripts/aeris.cron
MAILTO=""
1 * * * * cd weather;./aeris-weather-hourly
